# wearable/src/connman.service.in
[Unit]
Description=Connection service
After=tizen-system.target shutdown.target

[Service]
Type=dbus
BusName=net.connman
RemainAfterExit=yes
ExecStartPre=/usr/bin/dbus-send --system --dest=net.netconfig / net.netconfig.auto.activate
ExecStart=/usr/sbin/connmand

[Install]
WantedBy=multi-user.target
